# systemd/holler.service
[Unit]
Description=Holler MCP Gateway

[Service]
Type=simple
WorkingDirectory=%h/src/halfremembered-mcp
ExecStart=%h/src/halfremembered-mcp/target/release/holler serve --port 8080 --hootenanny tcp://localhost:5580 --hootenanny-pub tcp://localhost:5581
Restart=always
RestartSec=5

Environment=RUST_LOG=info
Environment=OTEL_EXPORTER_OTLP_ENDPOINT=localhost:4317

[Install]
WantedBy=default.target
