# systemd/holler.service
[Unit]
Description=Holler MCP Gateway
Documentation=https://github.com/halfremembered/hootenanny

[Service]
Type=simple
WorkingDirectory=%h/src/hootenanny
ExecStart=%h/src/hootenanny/target/release/holler serve
Restart=always
RestartSec=5
TimeoutStopSec=30

Environment=RUST_LOG=info

[Install]
WantedBy=default.target
