# systemd/luanette.service
[Unit]
Description=Luanette Lua Scripting Server
Documentation=https://github.com/halfremembered/hootenanny

[Service]
Type=simple
WorkingDirectory=%h/src/hootenanny
ExecStart=%h/src/hootenanny/target/release/luanette zmq --bind tcp://0.0.0.0:5570 --hootenanny tcp://localhost:5580
Restart=always
RestartSec=5
TimeoutStopSec=30

Environment=RUST_LOG=info

[Install]
WantedBy=default.target
