# systemd/vibeweaver.service
[Unit]
Description=Vibeweaver Python Kernel for AI Music Agents
Documentation=https://github.com/halfremembered/hootenanny
After=hootenanny.service
Wants=hootenanny.service

[Service]
Type=simple
WorkingDirectory=%h/src/hootenanny
ExecStart=%h/src/hootenanny/target/release/vibeweaver
Restart=always
RestartSec=5
TimeoutStopSec=30

Environment=RUST_LOG=info
Environment=PYTHONPATH=%h/.local/lib/python3.13/site-packages

[Install]
WantedBy=default.target
